# contrib/vault-watch.service
[Unit]
Description=VaultWatch — SAS HDD & RAID Monitor
After=network.target

[Service]
Type=simple
ExecStart=__BINARY_PATH__
Restart=on-failure
RestartSec=5s
# Run as root so smartctl can access raw disk devices without sudo
User=root
StandardInput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
