# contrib/vault-watch.openrc
#!/sbin/openrc-run

name="vault-watch"
description="VaultWatch — SAS HDD & RAID Monitor"
command="__BINARY_PATH__"
command_user="root"
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true

depend() {
    need net
}
